# crontab
*/15 * * * * speedtest-cli --secure --csv
